# gemini.service
[Unit]
Description=Gemini Polymarket Signal Engine
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=gemini
Group=gemini
WorkingDirectory=/opt/gemini
EnvironmentFile=/etc/gemini/gemini.env
ExecStartPre=/usr/bin/test -f /opt/gemini/backend.py
ExecStart=/usr/bin/python3 /opt/gemini/backend.py
Restart=always
RestartSec=10
KillSignal=SIGINT
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
